# server_config/supervisord.conf.sample
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/local/nginx/sbin/nginx

[program:unblockneteasemusic]
command=/usr/src/bin/unblockneteasemusic -p 8123 -k -q -r
